# docker/etc/supervisor/conf.d/noVNC.conf
[program:noVNC]
command=bash -l -c "./utils/launch.sh"
environment=HOME=/home/user,USER=user
directory=/opt/noVNC
autorestart=true
stdout_logfile=/var/log/supervisor/noVNC.log
redirect_stderr=true
